# ansible/roles/kubelish/templates/kubelish.service.j2
# roles/kubelish/templates/kubelish.service.j2
[Unit]
Description=Kubelish mDNS announcer for Kubernetes LoadBalancer services
After=avahi-daemon.service

[Service]
Type=notify
Environment="KUBECONFIG={{ kubeconfig_path }}"
ExecStart=/usr/local/bin/kubelish --namespace default watch
Restart=always
User=nobody
Group=nogroup

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
